# deploy/bbb26-fetch.timer
[Unit]
Description=BBB26 data fetch every 15 minutes

[Timer]
# Fire every 15 minutes, aligned to wall clock (:00, :15, :30, :45)
OnCalendar=*:0/15
# Small random delay to avoid thundering herd (0-30 seconds)
RandomizedDelaySec=30
# Run immediately if a trigger was missed (e.g., after reboot)
Persistent=true

[Install]
WantedBy=timers.target
