# deploy/bbb26-fetch.service
[Unit]
Description=BBB26 data fetch + build + push
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=bbb26
WorkingDirectory=/home/bbb26/BBB26
ExecStart=/usr/bin/python3 scripts/schedule_data_fetch.py --once --run-now --build --trigger-deploy
TimeoutStartSec=600
StandardOutput=journal
StandardError=journal

# Ensure git can find SSH key
Environment=HOME=/home/bbb26
Environment=GIT_SSH_COMMAND=ssh -i /home/bbb26/.ssh/bbb26_deploy -o StrictHostKeyChecking=accept-new
